# Codes/fishe_manager/startup_service_fishe.service
[Unit]
Description=Run fishe_manager

[Service]
Type=simple
ExecStart=/home/pi/Codes/fishe_manager/runfishe.sh

[Install]
WantedBy=default.target